# target/supervisor/conf.d/change-password.conf
[program:change-password]
command=change-password-server.py
autostart=true
autorestart=true
startretries=2
startsecs=1
redirect_stderr=true
user=change-password
killasgroup=true
stopasgroup=true